# systemd/aslan-health-check.timer
[Unit]
Description=Run Aslan Drive Health Check Daily at 8:00 AM
Requires=aslan-health-check.service

[Timer]
# Run daily at 8:00 AM (2 hours after data ingestion)
OnCalendar=*-*-* 08:00:00

# Randomize start time by up to 2 minutes
RandomizedDelaySec=120

# Ensure timer is persistent across reboots
Persistent=true

# If system was down during scheduled time, run immediately on boot
WakeSystem=false

[Install]
WantedBy=timers.target