# systemd/aslan-data-ingestion.timer
[Unit]
Description=Run Aslan Drive Data Ingestion Daily at 6:00 AM
Requires=aslan-data-ingestion.service

[Timer]
# Run daily at 6:00 AM
OnCalendar=*-*-* 06:00:00

# Randomize start time by up to 5 minutes to avoid resource contention
RandomizedDelaySec=300

# Ensure timer is persistent across reboots
Persistent=true

# If system was down during scheduled time, run immediately on boot
WakeSystem=false

[Install]
WantedBy=timers.target