# systemd/aslan-health-check.service
[Unit]
Description=Aslan Drive Health Check Morning Job
Requires=docker.service
After=docker.service
After=aslan-data-ingestion.service
Wants=aslan-health-check.timer

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/aslan-drive
Environment=DATABASE_URL=postgresql://trader:your_password@aslan_postgres:5432/aslan_drive
Environment=SLACK_WEBHOOK_URL=https://hooks.slack.com/services/YOUR/WEBHOOK/URL
Environment=LOG_LEVEL=info
ExecStartPre=/usr/bin/docker pull aslan_drive_health_check:latest
ExecStart=/opt/aslan-drive/scripts/run_health_check.sh
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/aslan
ReadWritePaths=/opt/aslan-drive

# Resource limits
MemoryLimit=256M
CPUQuota=25%

# Restart policy
RestartSec=2min
Restart=on-failure
RestartPreventExitStatus=2

[Install]
WantedBy=multi-user.target